# systemd/prometheus-pvdisplay-exporter.service
[Unit]
Description=Prometheus pvdisplay Exporter

[Service]
ExecStart=/usr/sbin/prometheus-pvdisplay-exporter
User=prometheus
Restart=always

[Install]
WantedBy=multi-user.target
